# dist/lib/systemd/system/siprec-drachtio.service
[Unit]
Description="SIPREC drachtio app"

[Service]
ExecStart=/usr/bin/node app.js siprec-drachtio
WorkingDirectory=__HOMEPATH__
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=SIPRECdrachtio
Environment=NODE_ENV=production PORT=8080

[Install]
WantedBy=multi-user.target
